# overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Auto-login root on the built-in display (tty1) for the test image — so a missing
# USB keyboard doesn't block getting a shell on-screen. CHANGE for any real build.
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin root --noclear %I 38400 linux
