# overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Auto-login the non-root 'deck' session user on the built-in display (tty1). deck's
# ~/.bash_profile execs the gamescope Big-Picture session; PipeWire (audio) only runs for a
# non-root user, which is why this is 'deck', not root. (root still reachable via SSH / other VTs.)
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin deck --noclear %I 38400 linux
