# overlay/etc/systemd/system/pocknix-gamescope-rt.service
[Unit]
Description=Elevate gamescope compositor threads to realtime (parity with ROCKNIX root session)
# Start early; do NOT order after multi-user.target — that pulls in the graphical Steam
# session bring-up and can settle late, leaving the start job queued forever.
After=basic.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pocknix-gamescope-rt
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
