# overlay/etc/systemd/system/pocknix-usbgadget.service
[Unit]
Description=pocknix USB-C network gadget (ssh root@10.66.0.1)
After=systemd-modules-load.service sys-kernel-config.mount
Wants=sys-kernel-config.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pocknix-usbgadget
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
