# overlay/etc/systemd/system/pocknix-lavd.service
[Unit]
Description=pocknix LAVD scheduler (scx_lavd, big.LITTLE / latency-aware, autopilot mode)
Documentation=https://github.com/sched-ext/scx
# Start early once the base system is up — do NOT order after multi-user.target, which
# pulls in the graphical Steam session bring-up and can settle late (or stall on the temp
# gamescope-diag), leaving the start job queued forever. The sched_ext/BTF guards below are
# satisfied from kernel init, so a scheduler needs no later ordering. WantedBy stays
# multi-user.target so it's still enabled at boot; default deps order it after basic.target.
After=basic.target
# Guards: skip cleanly (don't crash-loop) until the sched_ext + BTF kernel is deployed.
#   /sys/kernel/sched_ext  exists iff CONFIG_SCHED_CLASS_EXT=y in the running kernel.
#   /sys/kernel/btf/vmlinux exists iff DEBUG_INFO_BTF=y + pahole ran at build time.
ConditionPathExists=/sys/kernel/sched_ext
ConditionPathExists=/sys/kernel/btf/vmlinux
# If the BPF object can't load (e.g. a kernel missing the tracing/BPF-events stack that
# scx_lavd's bpf_trace_printk + futex/execve tracepoints need), give up after a few tries
# instead of crash-looping and flooding the journal.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
# The mode is NOT hardcoded here: pocknix-lavd-mode run reads the persisted choice from
# /var/lib/pocknix/lavd-mode (DEFAULT: --autopilot) and execs scx_lavd with it. Toggle with
# `pocknix-lavd-mode toggle` (or set autopilot/performance explicitly) — it rewrites the
# file and restarts this unit.
#
# --autopilot (default) = LAVD picks its own power mode from measured load:
#   required capacity <= ~half a little core  -> powersave  (littles-first, core compaction)
#   moderate (roughly 5-70% of total capacity) -> balanced   (compaction, energy-model order)
#   high (>~70%)                               -> performance (all cores, fast-cores-first)
# Chosen over --performance (2026-07-02) for battery/heat/fan noise: at idle and in menus the
# device now compacts onto little cores and drops clocks instead of biasing everything high.
# Trade-off (earlier on-device A/B): GPU-bound games often sit at moderate CPU util, so
# autopilot stays in *balanced* in-game and raw FPS measured a bit lower than --performance —
# flip to performance for a demanding title: pocknix-lavd-mode performance. Pairs with the
# schedutil governor (kernel default) so LAVD's per-core frequency requests
# (scx_bpf_cpuperf_set) actually take effect.
ExecStart=/usr/local/bin/pocknix-lavd-mode run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
