# overlay/etc/systemd/system/pocknix-volumed.service
[Unit]
Description=pocknix volume rocker -> PipeWire volume
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pocknix-volumed
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
