# overlay/etc/systemd/system/pocknix-diag.service
[Unit]
Description=pocknix first-boot diagnostic dump (writes pocknix-diag.txt to the boot partition)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
# delay so device probing AND a wifi association attempt have time to settle
ExecStartPre=/bin/sh -c 'sleep 30'
ExecStart=/usr/local/bin/pocknix-diag
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
