# overlay/etc/systemd/system/pocknix-fex-binfmt-off.service
# Disable the FEX x86/x86_64 binfmt_misc handlers at boot.
#
# WHY: the FEX package ships binfmt.d drop-ins that register FEX as the interpreter for x86 ELFs.
# For our use (Proton x86 *games*), that's not needed — the x86 Windows code goes Wine->FEX directly,
# not through Linux binfmt — and leaving it enabled breaks Steam's Proton compat-tool setup (the game
# fails to launch). ROCKNIX toggles binfmt off around the Steam session; we just keep it off.
#
# WHY A ROOT SERVICE (not in the launcher): writing /proc/sys/fs/binfmt_misc needs root. The
# boot-to-Steam session runs as the non-root 'deck' user, which CAN'T write it — so this must happen
# here, as root, after systemd-binfmt registers them. (This is the "binfmt-disable needs root" note.)
[Unit]
Description=pocknix: disable FEX x86 binfmt (breaks Proton compat-tool setup)
After=systemd-binfmt.service
Wants=systemd-binfmt.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'echo 0 > /proc/sys/fs/binfmt_misc/FEX-x86_64 2>/dev/null || true; echo 0 > /proc/sys/fs/binfmt_misc/FEX-x86 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
