# overlay/etc/systemd/system/pocknix-rt-demote.service
[Unit]
Description=pocknix — keep the Proton/FEX game off SCHED_RR (fixes gamescope GPU->0% hitching)
Documentation=file:///usr/share/doc/pocknix/gamescope-fps-investigation.md
# Wine RR-promotes the whole game (deck has rtprio 98); with the RT throttle off that causes
# priority-inversion stalls. This watcher demotes it back to SCHED_OTHER. See the script header.
After=basic.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pocknix-rt-demote-watch
# It's a cheap poll loop; keep it alive across game relaunches / session switches.
Restart=always
RestartSec=1
Nice=-5

[Install]
WantedBy=multi-user.target
