# overlay/etc/systemd/system/pocknix-lavd.service
[Unit]
Description=pocknix LAVD scheduler (scx_lavd, big.LITTLE / latency-aware, autopilot)
Documentation=https://github.com/sched-ext/scx
# Start early once the base system is up — do NOT order after multi-user.target, which
# pulls in the graphical Steam session bring-up and can settle late (or stall on the temp
# gamescope-diag), leaving the start job queued forever. The sched_ext/BTF guards below are
# satisfied from kernel init, so a scheduler needs no later ordering. WantedBy stays
# multi-user.target so it's still enabled at boot; default deps order it after basic.target.
After=basic.target
# Guards: skip cleanly (don't crash-loop) until the sched_ext + BTF kernel is deployed.
#   /sys/kernel/sched_ext  exists iff CONFIG_SCHED_CLASS_EXT=y in the running kernel.
#   /sys/kernel/btf/vmlinux exists iff DEBUG_INFO_BTF=y + pahole ran at build time.
ConditionPathExists=/sys/kernel/sched_ext
ConditionPathExists=/sys/kernel/btf/vmlinux
# If the BPF object can't load (e.g. a kernel missing the tracing/BPF-events stack that
# scx_lavd's bpf_trace_printk + futex/execve tracepoints need), give up after a few tries
# instead of crash-looping and flooding the journal.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
# --autopilot = load-adaptive power mode (balanced): eases clocks in menus/idle, ramps
# under game load, while always doing LAVD's latency-aware big.LITTLE placement. Pairs
# with the schedutil governor (kernel default) so LAVD drives per-core frequency too.
# Swap to --performance here if autopilot ever eases off at the wrong moment mid-game.
ExecStart=/usr/bin/scx_lavd --autopilot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
