# overlay/etc/systemd/system/pocknix-powerd.service
[Unit]
Description=pocknix power button -> short press suspends (long press stays logind poweroff)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pocknix-powerd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
