# packages/pocknix-desktop/pocknix-flathub.service
[Unit]
Description=pocknix: register the Flathub remote for Discover/flatpak
# Adding the remote fetches flathub.flatpakrepo (the GPG key + URL), so it needs the network the
# first time. After that the remote exists and ExecStart is a no-op even offline.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pocknix-flathub-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
