# devices/sm8250/packages/pocknix-bsp-sm8250/pocknix-sm8250-board-select.service
[Unit]
Description=pocknix sm8250 GRUB board-select correction
# Needs /flash (fstab, nofail) for grubenv; runs once per boot, exits fast.
After=flash.mount
ConditionPathIsMountPoint=/flash

[Service]
Type=oneshot
ExecStart=/usr/lib/pocknix/pocknix-sm8250-board-select

[Install]
WantedBy=multi-user.target
